# root/etc/s6-overlay/s6-rc.d/svc-calibre-web-automated/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Entry announcement (issue #868): every s6 unit names itself before doing
# work, so a `docker logs --timestamps` trace can attribute boot time to a
# service. This is a longrun, so the line also repeats on each supervisor
# restart, which is what makes a crash-loop visible.
echo "[svc-calibre-web-automated] Starting Calibre-Web server on port ${CWA_PORT_OVERRIDE:-8083}..."

export CALIBRE_DBPATH=/config

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${CWA_PORT_OVERRIDE:-8083}" \
        cd /app/calibre-web-automated s6-setuidgid abc python3 /app/calibre-web-automated/cps.py
