# root/etc/s6-overlay/s6-rc.d/metadata-change-detector/run
#!/usr/bin/with-contenv bash

echo "========== STARTING METADATA CHANGE DETECTOR ==========="

# Folder to monitor
WATCH_FOLDER="/app/calibre-web-automated/metadata_change_logs"
echo "[metadata-change-detector] Watching folder: $WATCH_FOLDER"

# Create the folder if it doesn't exist
install -d -o abc -g abc "$WATCH_FOLDER"

# Both watcher backends feed change-log *filenames* (one per line) into the
# debouncing dispatcher, which coalesces the burst of duplicate events a single
# save produces into at most one enforcement pass per file and then runs the
# enforcer. This unifies the inotify and polling paths behind one dedup so one
# save == one enforcement pass, with no "not found after 3 attempts" spam
# (fork #802). See scripts/metadata_change_dispatch.py.
DISPATCHER="/app/calibre-web-automated/scripts/metadata_change_dispatch.py"

# Monitor the folder for new files; on inotify errors, fall back to polling
run_fallback() {
        echo "[metadata-change-detector] Falling back to polling watcher (inotify unavailable or out of watches)" >&2
        python3 /app/calibre-web-automated/scripts/watch_fallback.py \
                --path "$WATCH_FOLDER" \
                --interval 5 \
                --exts "json,log" |
        while read -r events filepath; do
                basename -- "$filepath"
        done |
        python3 "$DISPATCHER" --watch-folder "$WATCH_FOLDER"
}

# Detect if running under Docker Desktop (Windows/macOS) and prefer polling
is_docker_desktop() {
        local osr mounts
        osr=$(cat /proc/sys/kernel/osrelease 2>/dev/null || true)
        if echo "$osr" | grep -qi 'microsoft'; then
                return 0
        fi
        if echo "$osr" | grep -qi 'linuxkit'; then
                return 0
        fi
        mounts=$(cat /proc/self/mountinfo 2>/dev/null || true)
        if echo "$mounts" | grep -Eqi '/host_mnt/|/Users/|osxfs|virtiofs.*docker'; then
                return 0
        fi
        return 1
}

# Prefer polling when running over network shares to better handle NFS/SMB semantics
if [ "${NETWORK_SHARE_MODE,,}" = "true" ] || [ "${NETWORK_SHARE_MODE}" = "1" ] || [ "${NETWORK_SHARE_MODE,,}" = "yes" ] || [ "${NETWORK_SHARE_MODE,,}" = "on" ]; then
        echo "[metadata-change-detector] NETWORK_SHARE_MODE=true detected; using polling watcher instead of inotify"
        run_fallback
        exit 0
fi

# Allow admins to force polling mode explicitly (defaults to inotify)
if [ "${CWA_WATCH_MODE:-inotify}" = "poll" ]; then
        run_fallback
        exit 0
fi

# Prefer polling on Docker Desktop (Windows/macOS) for reliable detection on host-mounted paths
if is_docker_desktop; then
        echo "[metadata-change-detector] Docker Desktop environment detected; using polling watcher instead of inotify"
        run_fallback
        exit 0
fi

(
        set -o pipefail
        s6-setuidgid abc inotifywait -m -e close_write -e moved_to --format '%f' --exclude '^.*\.(swp)$' "$WATCH_FOLDER" |
        python3 "$DISPATCHER" --watch-folder "$WATCH_FOLDER"
) || run_fallback